# bin/phpsock.service
[Unit]
Description=Phpsock library server daemon

[Service]
ExecStart=/usr/bin/node broker.js
WorkingDirectory=/home/phpsock

# try to restart service 5 times with 20s delay, then fail
Restart=always
RestartSec=20
StartLimitBurst=5
StartLimitInterval=150

Environment=ENV=production
User=nginx
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=phpsock

[Install]
WantedBy=multi-user.target
